# ai.openclaw.notification.high.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.openclaw.notification.high</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/amberives/.openclaw/workspace/scripts/flush-high-batch</string>
    </array>
    
    <key>StartCalendarInterval</key>
    <dict>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    
    <key>RunAtLoad</key>
    <false/>
    
    <key>KeepAlive</key>
    <false/>
    
    <key>StandardOutPath</key>
    <string>/Users/amberives/.openclaw/workspace/logs/notification-high.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/amberives/.openclaw/workspace/logs/notification-high.error</string>
</dict>
</plist>